# .config/systemd/user/swww-cycle.service
[Unit]
Description=Change swww-daemon wallpaper between day and night

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/swww.sh
